# aureox-backuper/crond.conf
[program:crond]
command=sh -c "/usr/sbin/crond -f -l 0"
autostart=true
autorestart=true
